# crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command

0 3 * * * root certbot renew >/proc/1/fd/1 2>/proc/1/fd/2
